# docker/config/cron.d/laravel-schedule
* * * * * su -s /bin/bash -c "cd /data/www && php artisan schedule:run >> /data/logs/cron.log 2>&1" www-data
